# aihelpers/shared/launchd/com.liftosaur.sidecar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.liftosaur.sidecar</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/anton/projects/liftosaur/aihelpers/shared/sidecar-runner.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/anton/projects/liftosaur</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/anton/projects/liftosaur/logs/sidecar-server/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/anton/projects/liftosaur/logs/sidecar-server/launchd-stderr.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
