# scripts/com.liftosaur.rollbar-orchestrator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.liftosaur.rollbar-orchestrator</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/anton/projects/liftosaur/scripts/rollbar-orchestrator-runner.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/anton/projects/liftosaur</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>1</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/anton/projects/liftosaur/logs/rollbar-orchestrator/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/anton/projects/liftosaur/logs/rollbar-orchestrator/launchd-stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
